# linux/lighttpd/examples/standalone/s6-rc.d/renew/run
#!/command/execlineb -P
with-contenv
/usr/local/sbin/renew
