# linux/lighttpd/examples/haproxy/s6-rc.d/haproxy/run
#!/command/execlineb -P
with-contenv
/usr/local/sbin/haproxy -W -db -f /usr/local/etc/haproxy/haproxy.cfg
