# linux/lighttpd/8.4/trixie/s6-rc.d/lighttpd/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
/usr/local/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf
