# linux/lighttpd/examples/haproxy/s6-rc.d/dumpssl/up
echo "Dumping SSL certificate on exit enabled"
